# crontab
*/5 * * * * XDG_RUNTIME_DIR=/run/user/$(id -u) ~/.local/share/scripts/battery-check 
